# app/flask_app.service
[Unit]
Description=Flask app service
After=network.target

[Service]
WorkingDirectory=/home/pavelmazikin/app
ExecStart=/home/pavelmazikin/app/venv/bin/python run.py

[Install]
WantedBy=multi-user.target